# stapler-scripts/litellm-proxy/com.litellm.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.litellm.proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/docker</string>
        <string>compose</string>
        <string>-f</string>
        <string>$HOME/dotfiles/stapler-scripts/litellm-proxy/docker-compose.yml</string>
        <string>up</string>
        <string>-d</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>$HOME/dotfiles/stapler-scripts/litellm-proxy</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/litellm-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/litellm-proxy.error.log</string>

    <key>StartInterval</key>
    <integer>0</integer>

    <!-- Wait for network to be available -->
    <key>NetworkState</key>
    <true/>
</dict>
</plist>
