# stapler-scripts/claude-proxy-rs/com.claude-proxy-rs-canary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Canary plist: Rust proxy on port 47001, Python proxy stays on 47000.
     Load with: launchctl load ~/Library/LaunchAgents/com.claude-proxy-rs-canary.plist
     Point a test Claude Code instance at it:
       alias test-claude='ANTHROPIC_BASE_URL=http://localhost:47001 claude'
     Unload with: launchctl unload ~/Library/LaunchAgents/com.claude-proxy-rs-canary.plist -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-proxy-rs-canary</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/claude-proxy-rs</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_CODE_OAUTH_TOKEN</key>
        <string>INSERT_OAUTH_TOKEN_HERE</string>
        <key>AWS_PROFILE</key>
        <string>Sandbox.AdministratorAccess</string>
        <key>AWS_REGION</key>
        <string>us-west-2</string>
        <key>PROXY_PORT</key>
        <string>47001</string>
        <key>COOLDOWN_SECONDS</key>
        <string>300</string>
        <key>REQUEST_TIMEOUT</key>
        <string>300</string>
        <key>BEDROCK_MAX_RETRIES</key>
        <string>3</string>
        <key>STAPLER_COMPRESS</key>
        <string>true</string>
        <key>CACHE_ALIGNER</key>
        <string>true</string>
        <key>VERBOSITY_LEVEL</key>
        <string>2</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/claude-proxy-rs-canary.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-proxy-rs-canary.error.log</string>
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
